# home/.config/conky/conky.service
[Unit]
Description=Run conky daemon to display system info on desktop

[Service]
ExecStart=/usr/bin/conky

